# etc/cron.d/aws-slsm
MAILTO=your-email@your-provider.com
* * * * * root /usr/local/bin/aws-slsm | logger -t aws-slsm
